# templates/xray-restart.service.j2
[Unit]
Description=Restart Xray VPN Service

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart xray.service
StandardOutput=journal
StandardError=journal
